# bot.conf
[program:telegrambot]
directory=/home/ubuntu/StartUp
command=python3 app.py

autostart=true
autorestart=true
enviroment=HOME="/home/ubuntu", USER="ubuntu"

stderr_logfile=/home/ubuntu/StartUp/logfile_err.log
stdout_logfile=/home/ubuntu/StartUp/logfile.log